# systemd_templates/dorb.target
[Unit]
Description=Systemd Master Starter
Wants=nginx.service postgres.service vpn.service redis_docker.service hello.service